# server/plaso-wkds.conf
start on runlevel [3]
respawn
exec /plaso/nodejs/bin/node /plaso/node-app/balance/balance.js 1>>/var/log/plaso/balance.log 2>>/var/log/plaso/balance.err.log
